# tachikoma-webserver-docker/src/assets/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:tachikoma_webserver]
command=/opt/tachikoma-webserver/bin/tachikoma-webserver
