# tachikoma-webserver-docker/src/assets/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info

[program:tachikoma_webserver]
command=/opt/tachikoma-webserver/bin/tachikoma-webserver

[program:rsyslogd]
command=/usr/sbin/rsyslogd -n