# tachikoma-postfix-docker/src/assets/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info

[program:postfix]
command=/opt/postfix.sh

# TODO Enable when working
#[program:opendkim]
#command=/usr/sbin/opendkim -f

# TODO disable when debugging
[program:tachikoma_utils]
command=/opt/tachikoma-postfix-utils/bin/tachikoma-postfix-utils
user=postfix

#[program:tachikoma_dev_proxy]
#command=/opt/proxy.py -L 2525 localhost 25
#user=postfix

[program:syslogd]
command=/usr/sbin/syslogd -n

[program:console_logger]
command=bash -c "tail -f /var/log/tachikoma/postfix-utils.log /var/log/mail.log > /dev/console"