# tachikoma-postfix-docker/src/assets/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:postfix]
command=/opt/postfix.sh

[program:opendkim]
command=/usr/sbin/opendkim -f
